# vnc-server/secure/supervisor.xvfb.conf
; Supervisor config file for Xvfb

[program:Xvfb]
command=/usr/bin/xvfb-run -n 1 -s "+xinerama +extension GLX +extension RANDR -screen 0 1024x780x24" -w 1 /usr/bin/openbox-session
priority=10
stdout_logfile=/var/log/Xvfb.log
stderr_logfile=/var/log/Xvfb-error.log
