# vnc-server/secure/supervisor.sshd.conf
; Supervisor config file for SSH server

[program:sshd]
command=/usr/sbin/sshd -D -e
priority=30
stdout_logfile=/var/log/sshd.log
redirect_stderr=true
